# deploy/claude-rpg.service
[Unit]
Description=ClaudeRPG game server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__REPO__
EnvironmentFile=/etc/claude-rpg.env
ExecStart=__REPO__/scripts/pi/run-server.sh
Restart=always
RestartSec=3
# Give the DB a moment to flush on stop
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
